# conf/hkradioweb.conf
[program:hkradioweb]
directory = /mnt/raspi/radio-control
command=python /mnt/raspi/radio-control/code.py

