# conf/hkradio.conf
[program:hkradio]
command=python /mnt/raspi/webstore/pyServ/uRTHK.py

